# systemd/iceberg-maintenance.timer
[Unit]
Description=Weekly Iceberg maintenance timer

[Timer]
OnCalendar=Sun 03:30 UTC
Persistent=true

[Install]
WantedBy=timers.target
